# supervisord.conf
[program:pinpin]
user=root
command=gunicorn -c gun.conf myapp:app
autostart=true
autorestart=true
stderr_logfile=err.log
stdout_logfile=out.log
stopsignal=INT

[supervisord]
